# notifications.service.in
[Unit]
Description=Notifications daemon
Requires=vconf-setup.service
After=vconf-setup.service

[Service]
ExecStart=/usr/bin/notification-service

[Install]
WantedBy=default.target
