# notifications.service.in
[Unit]
Description=Notifications daemon

[Service]
ExecStart=/usr/bin/notification-service

[Install]
WantedBy=graphical.target
