# notifications-display-x11.service.in
[Unit]
Description=Notifications Popup daemon
Requires=notifications.service display-manager-run.service
After=notifications.service display-manager-run.service

[Service]
ExecStart=/usr/bin/notification-display-service
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
