# notifications-display.service.in
[Unit]
Description=Notifications Popup daemon
Requires=notifications.service display-manager-run.service
After=notifications.service display-manager-run.service

[Service]
ExecStart=/usr/bin/notification-display-service
EnvironmentFile=/etc/sysconfig/weston

[Install]
WantedBy=graphical.target
