# notifications.service.in
[Service]
ExecStart=/usr/bin/notification-service
