# notifications-display-ivi.service.in
[Unit]
Description=Notifications Popup daemon
Requires=notifications.service
After=notifications.service

[Service]
ExecStart=/usr/bin/notification-display-service
Environment=XDG_RUNTIME_DIR=/run/user/5000

[Install]
WantedBy=default.target
