# notifications_display.service.in
[Unit]
Description=Notifications Popup daemon

[Service]
Type=simple
Environment=XDG_RUNTIME_DIR=/run/display
ExecStart=/usr/bin/notification-display-service

