# notifications.socket.in
[Socket]
ListenStream=/tmp/.notification.service

[Install]
WantedBy=sockets.target
